# plists/launchd.meshcore-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>launchd.meshcore-bridge</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/muehlbucks/.opentakserver_venv/bin/python</string>
      <string>/Users/muehlbucks/ots/meshcore_bridge.py</string>
      <string>--config</string>
      <string>/Users/muehlbucks/ots/config.yml</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/muehlbucks/ots</string>
    <key>UserName</key>
    <string>muehlbucks</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
      <key>Crashed</key>
      <true/>
      <key>OtherJobEnabled</key>
      <dict>
        <key>homebrew.mxcl.rabbitmq</key>
        <true/>
      </dict>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/muehlbucks/ots/logs/meshcore-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/muehlbucks/ots/logs/meshcore-bridge.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/muehlbucks/.opentakserver_venv/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
