# plists/launchd.meshcore-mqtt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>launchd.meshcore-mqtt</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/meshcore-mqtt</string>
      <string>--config-file</string>
      <string>/Users/muehlbucks/.openclaw/workspace/projects/atak-nj/meshcore-mqtt/config.yaml</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/muehlbucks/.openclaw/workspace/projects/atak-nj/meshcore-mqtt</string>
    <key>UserName</key>
    <string>muehlbucks</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
      <key>Crashed</key>
      <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/muehlbucks/ots/logs/meshcore-mqtt.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/muehlbucks/ots/logs/meshcore-mqtt.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
