# pi/systemd/adsbcot.service
[Unit]
Description=ADS-B to CoT Gateway
After=readsb.service eud-handler-tcp.service
Wants=readsb.service

[Service]
Type=simple
User=__BIGMAC_USER__
WorkingDirectory=__BIGMAC_DATA__
ExecStart=__BIGMAC_VENV__/bin/adsbcot -c __BIGMAC_DATA__/adsbcot/config.ini
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
