# pi/systemd/eud-handler-tcp.service
[Unit]
Description=OTS EUD Handler (TCP port 8088)
After=opentakserver.service rabbitmq-server.service
Requires=rabbitmq-server.service

[Service]
Type=simple
User=__BIGMAC_USER__
WorkingDirectory=__BIGMAC_DATA__
ExecStart=__BIGMAC_VENV__/bin/eud_handler
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
