# pi/systemd/meshtastic-bridge.service
[Unit]
Description=Meshtastic Serial to OTS CoT Bridge
After=rabbitmq-server.service
Requires=rabbitmq-server.service

[Service]
Type=simple
User=__BIGMAC_USER__
WorkingDirectory=__BIGMAC_DATA__
ExecStart=__BIGMAC_VENV__/bin/python3 __BIGMAC_DATA__/meshtastic_bridge.py --config __BIGMAC_DATA__/config.yml
Restart=on-failure
RestartSec=15
# Give the USB serial device time to enumerate
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
